# root/etc/s6-overlay/s6-rc.d/init-aria2-config/up
/etc/s6-overlay/s6-rc.d/init-aria2-config/run